# recipes-support/blink-led/blink-led/blink-led.service
[Unit]
Description=LED-blink test
After=rauc.service

[Service]
Type=simple
ExecStart=/usr/bin/blink-led

[Install]
WantedBy=multi-user.target
