# recipes-support/blink-led/blink-led/blink-led.timer
[Unit]
Description=Periodically check RAUC status for blink-led service

[Timer]
OnActiveSec=3
OnUnitActiveSec=3

[Install]
WantedBy=timers.target
