# recipes-connectivity/awsclient/awsclient/awsclient.timer
[Unit]
Description=Timer for executing OSB AWS Client periodically

[Timer]
OnActiveSec=1m
OnUnitActiveSec=30m

[Install]
WantedBy=basic.target
